# dovecot/supervisord.ini
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
user = root

[program:rsyslog]
command = /usr/sbin/rsyslogd -n
priority = 10

[program:healthcheck_permissions]
startsecs = 0
exitcodes = 0
priority = 40
command = /bin/sh -c 'chmod 755 /health-check.sh'

[program:vmail_permissions]
startsecs = 0
exitcodes = 0
priority = 50
command = /bin/sh -c 'chown -R vmail:vmail /srv/vmail'

[program:dovecot]
command = /usr/sbin/dovecot -F
priority = 60

[program:cron]
command = /usr/sbin/cron -f
