# postfix/supervisord.ini
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
user = root

[program:rsyslog]
command = /usr/sbin/rsyslogd -n
priority = 10

[program:postfix_configs]
startsecs = 0
exitcodes = 0,1
command = /bin/sh -c '[ -z "`ls -1 /etc/postfix`" ] && cp -R /etc/_postfix/* /etc/postfix/'
priority = 20

[program:postfix_configs_permissions]
startsecs = 0
exitcodes = 0
priority = 30
command = /bin/sh -c 'chown -R root:root /etc/postfix && chmod -R 755 /etc/postfix'

[program:healthcheck_permissions]
startsecs = 0
exitcodes = 0
priority = 40
command = /bin/sh -c 'chmod 755 /health-check.sh'

[program:posfix]
command = /usr/sbin/postfix start-fg
priority = 50

[program:crond]
command = /usr/sbin/crond -f
